# novpn-c/debian/prerm

systemctl stop novpn-client
systemctl disable novpn-client



